# bgp/supervisord.conf
[supervisord]
nodaemon=true

[program:bgp]
command=/data/start-input-bgp.sh
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/bgp.log
stdout_logfile=/var/log/bgp.log
